# fleet-watchdog.service
[Unit]
Description=Fleet watchdog check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fleet-watchdog
