# fleet-watchdog.service
[Unit]
Description=Fleet watchdog check

[Service]
Type=oneshot
PassEnvironment=GH_TOKEN ZABBIX_SERVER ZABBIX_PORT ZABBIX_API_URL ZABBIX_API_TOKEN
ExecStart=/usr/local/bin/fleet-watchdog
